# services/page_rank.service
[Unit]
Description=Page rank service
After=multi-user.target
[Service]
Type=simple
Restart=on-failure
ExecStart=/opt/rh/rh-python38/root/usr/bin/python /var/www/html/search-engine/run_page_rank.py
[Install]
WantedBy=multi-user.target